# ansible/roles/vpn/templates/dcv-vpn-watch.timer.j2
[Unit]
Description=Run DCV VPN connect/disconnect watcher every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Unit=dcv-vpn-watch.service

[Install]
WantedBy=timers.target
